# launchd/com.priority-forge.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.priority-forge.watchdog</string>
    
    <key>Program</key>
    <string>/bin/bash</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            export PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin
            export HOME=/Users/nicholasfleischhauer
            LOG=/tmp/priority-forge-watchdog.log
            echo "[$(date)] Watchdog started" >> $LOG
            while true; do
                # Check backend
                if ! curl -s --connect-timeout 2 http://127.0.0.1:3456/version > /dev/null 2>&amp;1; then
                    echo "[$(date)] Backend down - restarting" >> $LOG
                    launchctl kickstart -k gui/$(id -u)/com.priority-forge.backend 2>/dev/null
                fi
                # Check frontend  
                if ! curl -s --connect-timeout 2 http://127.0.0.1:5173 > /dev/null 2>&amp;1; then
                    echo "[$(date)] Frontend down - restarting" >> $LOG
                    launchctl kickstart -k gui/$(id -u)/com.priority-forge.frontend 2>/dev/null
                fi
                sleep 30
            done
        </string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>StandardOutPath</key>
    <string>/tmp/priority-forge-watchdog-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/priority-forge-watchdog-stderr.log</string>
</dict>
</plist>
