# launchd/com.priority-forge.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.priority-forge.backend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/nicholasfleischhauer/Desktop/priority-forge/dist/index.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/nicholasfleischhauer/Desktop/priority-forge</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/nicholasfleischhauer/Desktop/priority-forge/.pids/launchd-backend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/nicholasfleischhauer/Desktop/priority-forge/.pids/launchd-backend-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
